# systemd/cerebro2mqtt.service
[Unit]
Description=Cerebro2MQTT bridge service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/cerebro2mqtt
Environment=CEREBRO_CONFIG=/opt/cerebro2mqtt/config/config.json
Environment=LOG_LEVEL=INFO
ExecStart=/usr/bin/python3 -m app.main
ExecStartPre=/bin/sh -c 'systemctl stop serial-getty@ttyS0.service 2>/dev/null || true'
ExecStartPre=/bin/sh -c 'systemctl stop serial-getty@serial0.service 2>/dev/null || true'
ExecStartPre=/bin/sh -c 'fuser -k /dev/ttyS0 >/dev/null 2>&1 || true'
Restart=always
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target
